# script/health-ping.cron
* * * * * curl -s https://momemtum-server.onrender.com/api/health > /dev/null